# project2/software/build/linux/rootfs/packages-repo/targetroot/var/lib/dpkg/info/run-postinsts.list
/.
/usr
/usr/share
/usr/share/run-postinsts
/usr/share/run-postinsts/run-postinsts.awk
/etc
/etc/init.d
/etc/init.d/run-postinsts
